# playbook/roles/comedy_worker/files/comedy_worker.conf.j2
description "comedy_worker"

start on startup
stop on shutdown

setuid app
setgid app
chdir /app/current
env RAILS_ENV=production

respawn

exec /app/current/bin/comedy_worker
