# playbook/roles/puma/files/puma.conf.j2
description "puma"

start on startup
stop on shutdown

setuid app
setgid app
chdir /app/current

respawn

reload signal SIGUSR1

exec /app/current/bin/puma -C /app/current/config/puma.rb